# deploy/order-workflow-sample.service
[Unit]
Description=Order Workflow
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=order-workflow
Group=order-workflow
WorkingDirectory=/opt/order-workflow
Environment=PORT=3000
ExecStart=/usr/bin/python3 /opt/order-workflow/src/server.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
